# deploy/jarviddin.service.example
[Unit]
Description=Jarviddin Telegram orchestrator
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/jarviddin
EnvironmentFile=/opt/jarviddin/.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
